# security/vpn/systemd/aladdin-vpn.service
[Unit]
Description=ALADDIN VPN Service
Documentation=https://aladdin-vpn.com/docs
After=network.target network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=aladdin
Group=aladdin
WorkingDirectory=/opt/aladdin-vpn
ExecStart=/opt/aladdin-vpn/venv/bin/python /opt/aladdin-vpn/service_orchestrator.py
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aladdin-vpn

# Ограничения ресурсов
LimitNOFILE=65535
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aladdin-vpn/logs /opt/aladdin-vpn/data /opt/aladdin-vpn/backups
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Переменные окружения
Environment=PYTHONPATH=/opt/aladdin-vpn
Environment=ALADDIN_ENV=production
Environment=ALADDIN_LOG_LEVEL=INFO
Environment=ALADDIN_CONFIG_PATH=/opt/aladdin-vpn/config

[Install]
WantedBy=multi-user.target