# security/vpn/systemd/aladdin-vpn-backup.service
[Unit]
Description=ALADDIN VPN Backup Service
Documentation=https://aladdin-vpn.com/docs/backup
After=network.target aladdin-vpn.service
Wants=aladdin-vpn.service
Requires=network.target

[Service]
Type=simple
User=aladdin
Group=aladdin
WorkingDirectory=/opt/aladdin-vpn
ExecStart=/opt/aladdin-vpn/venv/bin/python /opt/aladdin-vpn/backup/backup_manager.py
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aladdin-vpn-backup

# Ограничения ресурсов
LimitNOFILE=8192
LimitNPROC=512
MemoryMax=256M
CPUQuota=25%

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aladdin-vpn/backups /opt/aladdin-vpn/logs /opt/aladdin-vpn/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Переменные окружения
Environment=PYTHONPATH=/opt/aladdin-vpn
Environment=ALADDIN_ENV=production
Environment=ALADDIN_LOG_LEVEL=INFO
Environment=BACKUP_DIRECTORY=/opt/aladdin-vpn/backups
Environment=MAX_BACKUPS=30
Environment=BACKUP_COMPRESSION_LEVEL=6

[Install]
WantedBy=multi-user.target