# security/vpn/systemd/aladdin-vpn-graphql.service
[Unit]
Description=ALADDIN VPN GraphQL API
Documentation=https://aladdin-vpn.com/docs/api
After=network.target aladdin-vpn.service
Wants=aladdin-vpn.service
Requires=network.target

[Service]
Type=simple
User=aladdin
Group=aladdin
WorkingDirectory=/opt/aladdin-vpn
ExecStart=/opt/aladdin-vpn/venv/bin/python /opt/aladdin-vpn/api/graphql_api.py
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aladdin-vpn-graphql

# Ограничения ресурсов
LimitNOFILE=32768
LimitNPROC=2048
MemoryMax=1G
CPUQuota=100%

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aladdin-vpn/logs /opt/aladdin-vpn/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Переменные окружения
Environment=PYTHONPATH=/opt/aladdin-vpn
Environment=ALADDIN_ENV=production
Environment=ALADDIN_LOG_LEVEL=INFO
Environment=GRAPHQL_HOST=0.0.0.0
Environment=GRAPHQL_PORT=8080

[Install]
WantedBy=multi-user.target